# etc/supervisor/conf.d/rails_app.conf
[program:rails_app]
directory=/opt/rails_app
command=bundle exec puma --port 80 --config ./config/puma.rb
