# etc/supervisor/conf.d/redis.conf
[program:redis]
command=redis-server /etc/redis.conf
