# etc/supervisor/conf.d/sidekiq.conf
[program:sidekiq]
directory=/opt/rails_app
command=bundle exec sidekiq -C config/sidekiq.yml
